# tools/ulogd2-clickhouse/ulogd2.service
[Unit]
Description=ulogd2 NFLOG Daemon — StaySuite SNI Capture Pipeline
Documentation=https://www.netfilter.org/projects/ulogd/
After=network-online.target firewalld.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/ulogd2/sbin/ulogd -c /usr/local/ulogd2/etc/ulogd.conf
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

# Pre-create log directories
ExecStartPre=/bin/mkdir -p /var/log/ulogd/json
ExecStartPre=/bin/chmod 755 /var/log/ulogd /var/log/ulogd/json

# Load kernel module if missing
ExecStartPre=-/sbin/modprobe nfnetlink_log

[Install]
WantedBy=multi-user.target
